# Binwalk_Extracted_Firmware/squashfs-root/etc_ro/inittab
::sysinit:/etc_ro/rcS
ttyS1::respawn:/sbin/getty 115200 ttyS1
